# systemd/todo-reminder.service
[Unit]
Description=Todo reminder notifications

[Service]
Type=oneshot
ExecStart=%h/.local/bin/todo notify
