# systemd/todo-reminder.timer
[Unit]
Description=Run todo reminders every minute

[Timer]
OnUnitActiveSec=60s
Persistent=true

[Install]
WantedBy=timers.target
